# templates/etc/profile.d/go.sh.j2
# {{ ansible_managed }}

export GOPATH="{{ go_install_prefix }}/go";
export PATH="${PATH}:${GOPATH}/bin";
